# docker/crontab
# ISP Status Page - Cron Configuration
# Run monitor checks every minute

* * * * * www-data /usr/local/bin/monitor-check-cron.sh >> /var/www/html/logs/cron.log 2>&1

# Run cleanup daily at 3 AM
0 3 * * * www-data /usr/local/bin/cleanup-cron.sh >> /var/www/html/logs/cleanup.log 2>&1

# Run backup daily at 2 AM
0 2 * * * www-data /usr/local/bin/backup-cron.sh >> /var/www/html/logs/backup.log 2>&1

# Empty line required at end of file
