# docker/crontab
# ISP Status Page - Cron Configuration
# Updated: Phase 4 cleanup — monitor_check and webhook_retry removed
# (scheduler daemon + queue workers handle those now)

# Safety fallback: run scheduler --once every minute
# In case the scheduler daemon container dies, this ensures checks still get pushed to the queue
* * * * * www-data cd /var/www/html && bin/cake scheduler --once >> /var/www/html/logs/cron.log 2>&1

# Run cleanup daily at 3 AM
0 3 * * * www-data /usr/local/bin/cleanup-cron.sh >> /var/www/html/logs/cleanup.log 2>&1

# Run backup daily at 2 AM
0 2 * * * www-data /usr/local/bin/backup-cron.sh >> /var/www/html/logs/backup.log 2>&1

# Empty line required at end of file
